# profiles/sysadmin/seed/etc/cron.d/backupjob
* * * * * backup /usr/local/bin/backup.sh
